# packaging/deb/DEBIAN/postrm
#!/bin/sh
set -e
if command -v systemctl >/dev/null 2>&1; then
  case "$1" in
    remove|purge)
      systemctl disable --now sunlionet-inside.service >/dev/null 2>&1 || true
      systemctl disable --now sunlionet-outside.service >/dev/null 2>&1 || true
      ;;
  esac
  systemctl daemon-reload || true
fi
